# scripts/com.benefron.newsdigest.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.benefron.newsdigest</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/automation/.venv/bin/python3</string>
        <string>__REPO_ROOT__/automation/run_daily.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__/automation</string>

    <!-- 07:00 / 10:00 / 12:00 / 15:00 / 17:00 / 19:00 / 21:00 daily.
         RunAtLoad is deliberately absent: launchd evaluates any
         StartCalendarInterval slot(s) missed while the Mac was asleep/off
         on its next wake, which is what gives the "runs at the closest
         available time" behavior instead of firing on every reload.
         run_daily.py's own MIN_HOURS_BETWEEN_RUNS guard absorbs any
         multiple-missed-slots catch-up so it doesn't double-run. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>10</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>12</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>19</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/automation/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/automation/logs/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__CLAUDE_BIN_DIR__:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
    </dict>
</dict>
</plist>
